# package-scripts/bixby/postinst
#!/bin/bash

# bixby post-install script

url="https://s3.bixby.io"

BIXBY_HOME=/opt/bixby

function bootstrap_repo() {
  if [ -d /opt/bixby/repo ]; then
    return 0;
  fi

  cd $BIXBY_HOME
  mkdir repo
  cd repo
  ver=`\\curl -sL $url/latest-repo`
  pkg="repo-$ver.tar.gz"
  echo "downloading $pkg"
  wget -q $url/repo/$pkg
  tar -xzf $pkg
  mv bixby-repo vendor
  rm -f $pkg

  return 0
}

function is_registered {
  [[ -f $BIXBY_HOME/etc/bixby.yml && -f $BIXBY_HOME/etc/server.pub ]]
}

function upgrade_bixby() {
  if [[ ! -f /etc/init.d/bixby ]]; then
    # install init script & god confs
    cd $BIXBY_HOME
    # change to bixby-agent gem dir so we can copy files to BIXBY_HOME
    cd $(readlink -f $(dirname $(embedded/bin/gem which bixby-agent))/..)

    cp -a etc/bixby-god.initd /etc/init.d/bixby
    chmod 755 /etc/init.d/bixby

    # install system service
    issue=`cat /etc/issue`
    [[ $issue =~ ^Ubuntu ]] && [[ `which update-rc.d` ]]  && update-rc.d bixby defaults
    [[ $issue =~ ^CentOS ]] && [[ `which chkconfig` ]]    && chkconfig --add bixby && chkconfig bixby on

    # setup god configs
    mkdir -p $BIXBY_HOME/etc
    cp -a etc/bixby.god etc/god.d $BIXBY_HOME/etc/
    cd $BIXBY_HOME

    # if is_registered; then
    #   /opt/bixby/bin/bixby-agent stop
    # fi
  fi

  # restart service using bixby
  # if is_registered; then
  #   /etc/init.d/bixby restart
  # fi
}

bootstrap_repo
chown -R bixby:bixby $BIXBY_HOME
upgrade_bixby
